# daemons/rclone.service
[Unit]
Description=Rclone sync with Google Drive
After=network-online.target
OnFailure=rclone-fail.service

[Service]
Type=oneshot
# Run the bisync through the helper script instead of calling rclone directly,
# so the bisync flags live in one place (config.sh, sourced by the script)
# rather than being duplicated in this unit. sync-service exits with rclone's
# status, so any non-zero exit is treated as failure and triggers
# OnFailure=rclone-fail.service above.
ExecStart=%h/bin/rclone4gdrive/rclone4gdrive sync-service
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
